# deploy/webfileserver.service.example
[Unit]
Description=WebFileServer
After=network.target mysql.service redis-server.service

[Service]
Type=simple
WorkingDirectory=/opt/WebFileServer
ExecStart=/opt/WebFileServer/main
Restart=on-failure
RestartSec=3
User=webfile
Group=webfile

[Install]
WantedBy=multi-user.target
